# data_acquisition/data_acquisition_stream/data/crontab
30 * * * * cd /app && python /app/main.py 1 >> /var/log/cron.log 2>&1